# slave/root/etc/services.d/sshd/run
#!/usr/bin/with-contenv sh

/usr/sbin/sshd -D
