# slave/root/etc/services.d/autofs/run
#!/usr/bin/with-contenv sh

/usr/sbin/automount --foreground --verbose /etc/auto.master
